# deploy/quorum-runner.service
[Unit]
Description=Quorum continuous debate runner
Documentation=https://github.com/furkankoykiran/Quorum
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/Quorum
ExecStart=/root/.pyenv/shims/uv run python -m apps.orchestrator.cli run --symbol SOL/USDT --interval 900
Restart=on-failure
RestartSec=30

Environment=PATH=/root/.pyenv/shims:/root/.local/bin:/root/.nvm/versions/node/v24.14.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/root

StandardOutput=append:/root/Quorum/data/runner.log
StandardError=append:/root/Quorum/data/runner.log

KillSignal=SIGTERM
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
